# crates/terraphim_llm_proxy/debian/postinst
#!/bin/bash

# Post-installation script for terraphim-llm-proxy

set -e

case "$1" in
    configure)
        # Create user and group if they don't exist
        if ! getent passwd terraphim-llm-proxy >/dev/null; then
            adduser --system --group --no-create-home --home /var/lib/terraphim-llm-proxy \
                --shell /usr/sbin/nologin terraphim-llm-proxy
        fi

        # Create data directory
        mkdir -p /var/lib/terraphim-llm-proxy
        chown terraphim-llm-proxy:terraphim-llm-proxy /var/lib/terraphim-llm-proxy

        # Enable systemd service
        if systemctl --version >/dev/null 2>&1; then
            systemctl daemon-reload
            systemctl enable terraphim-llm-proxy.service || true
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0