# systemd/adf-orchestrator.service.d/memory.conf
# systemd drop-in for adf-orchestrator: memory pressure threshold.
#
# Install on bigbox:
#   sudo mkdir -p /etc/systemd/system/adf-orchestrator.service.d
#   sudo install -m 644 systemd/adf-orchestrator.service.d/memory.conf \
#       /etc/systemd/system/adf-orchestrator.service.d/memory.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart adf-orchestrator
#
# Verify:
#   systemctl show adf-orchestrator | grep -i memory
#
# Rationale: prior to 2026-05-23 the orchestrator reached 84.3 G of its
# 90 G MemoryHigh ceiling after long uptime + leaked agent tasks. This
# drop-in lowers the soft pressure point to 80 G so the kernel applies
# back-pressure earlier and we get more headroom before MemoryMax kicks
# in. Combined with OnFailure restart unit, this gives graceful recovery.

[Service]
MemoryHigh=80G
